# scripts/net-auto-switch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.octopusgarage.net-auto-switch</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PROJECT_DIR__/scripts/launchd-wrapper.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_DIR__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <!-- Routine logs live in the daily-rotated ~/Library/Logs/net_auto_switch.log,
       so stdout is discarded to avoid an unbounded duplicate. stderr is kept as a
       small safety net for crashes that happen before logging is configured. -->
  <key>StandardOutPath</key>
  <string>/dev/null</string>

  <key>StandardErrorPath</key>
  <string>__PROJECT_DIR__/logs/launchd.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
